# recipes-warp7/modbus-rtu/files/modbus-rtu.service
[Unit]
Description=modbusRTUServer daemon service

[Service]
ExecStart=/usr/bin/modbusRTUServer

[Install]
WantedBy=multi-user.target
