# recipes-warp7/service-gateway/files/service-gateway.service
[Unit]
Description=ServiceGateway daemon service

[Service]
ExecStart=/usr/bin/ServiceGateway

[Install]
WantedBy=multi-user.target
