# recipes-warp7/clicks-board-init/clicks-board-init/mikrobus.service
[Unit]
Description=Click Board init daemon service

[Service]
ExecStart=/usr/bin/init-mikrobus.sh

[Install]
WantedBy=multi-user.target
